# backend/deploy/pussel-deploy.service
# Template for /etc/systemd/system/pussel-deploy.service — set User= and the
# ExecStart path to the deploy user and checkout on your server before
# installing. Triggered by pussel-deploy.timer; see backend/deploy/README.md.
[Unit]
Description=Deploy pussel backend when origin/main moves
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=DEPLOY_USER
ExecStart=/path/to/pussel/repo/backend/deploy/pussel-deploy.sh
