# backend/deploy/pussel-deploy.timer
# Template for /etc/systemd/system/pussel-deploy.timer.
[Unit]
Description=Poll origin/main every 5 minutes and deploy the pussel backend

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
