# backend/deploy/pussel-deploy.service
# Installed at /etc/systemd/system/pussel-deploy.service on delectosoft.
# Triggered by pussel-deploy.timer; see backend/deploy/README.md.
[Unit]
Description=Deploy pussel backend when origin/main moves
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=claus
ExecStart=/home/claus/pussel/repo/backend/deploy/pussel-deploy.sh
